# scripts/systemd/manta-server.service
# manta-server.service — systemd --user unit for the manta box server.
#
# Mirrors the hand-rolled unit documented in AGENTS.md ("Mobile / web client"):
# runs `node src/server/index.mjs` bound to loopback, restarted on failure,
# started at boot via `loginctl enable-linger $USER`.
#
# The installer (scripts/install.sh) substitutes the @@…@@ placeholders and
# writes the result to ~/.config/systemd/user/manta-server.service, then:
#   systemctl --user daemon-reload
#   systemctl --user enable --now manta-server
#
# Bind is 127.0.0.1 on purpose: the server is NOT exposed directly. The
# installer (scripts/install.sh) fronts it with Caddy on
# https://<box_id>.boxes.mantaui.com, which reverse-proxies loopback to
# 127.0.0.1:8787. Pairing codes are loopback-gated so `manta pair` must run
# on the box.

[Unit]
Description=manta box server (mobile/web + push fanout via gateway)
Documentation=https://github.com/antoinedc/MantaUI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@@MANTA_HOME@@
ExecStart=@@NODE_BIN@@ @@MANTA_HOME@@/src/server/index.mjs
Environment=MANTA_MOBILE_HOST=127.0.0.1
Environment=MANTA_MOBILE_PORT=@@MANTA_PORT@@
# Optional second listener bound to the Tailscale interface (BET-266).
# Set to a Tailscale IPv4 by install.sh on the tailnet ingress path;
# left empty (single-quoted empty string) on the public path so the
# server's tailnet listener is a no-op. WireGuard encrypts the tailnet
# hop, so the box is reachable WITHOUT TLS on this listener.
Environment=MANTA_TAILNET_HOST=@@MANTA_TAILNET_HOST@@
Restart=on-failure
RestartSec=2
# Give the server a moment to flush on stop before SIGKILL.
TimeoutStopSec=10

[Install]
WantedBy=default.target
